# deploy/systemd/paperclip.service
[Unit]
Description=Paperclip UI
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=paper
WorkingDirectory=/opt/paper-company
ExecStart=/opt/paper-company/.venv/bin/python scripts/paperclip_server.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
