# deploy/systemd/paper-company-daily.service
[Unit]
Description=Paper Company Morning Signal run

[Service]
Type=oneshot
User=paper
WorkingDirectory=/opt/paper-company
EnvironmentFile=/opt/paper-company/.env
ExecStart=/opt/paper-company/.venv/bin/python scripts/explore_daily.py
