# deploy/systemd/paper-company-daily.timer
[Unit]
Description=Paper Company Morning Signal daily trigger

[Timer]
OnCalendar=*-*-* 07:00:00 Asia/Seoul
Persistent=true

[Install]
WantedBy=timers.target
